# deployment/tiles.dev.publicamundi.eu/etc/supervisor/conf.d/mapproxy.conf
[program:mapproxy]
directory=/var/local/mapproxy/etc
command=/var/local/mapproxy/bin/gunicorn -k eventlet -w 4 -b :8000 config:application
autostart=true
