# machines/geoserver.dev.publicamundi.eu/etc/supervisor/conf.d/reloader.conf
[program:reloader]
directory=/root/reloader
command=python /root/reloader/reload.py
autostart=true
